# iac/containers/php/build/supervisor/queue-work.conf
[program:queue-work]
process_name=%(program_name)s_%(process_num)02d
directory=/var/www
command=php artisan queue:work
autostart=true
autorestart=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/queue-work.log